# conf/cron.daily
* * * * * root /usr/bin/reprepro -b /srv/buildd processincoming incoming
#* * * * * root /usr/bin/reprepro -b /srv/buildd processincoming lucid
#* * * * * root /usr/bin/reprepro -b /srv/buildd processincoming natty
